# testprj/crontab
* 10 * * * /usr/local/bin/python3 /code/manage.py notify

